# pi/services/bor-hanger.service
[Unit]
Description=BOR Systems WiFi hanger daemon
After=network-online.target
Wants=network-online.target
# Don't start the main daemon until the customer has finished the BLE Wi-Fi
# onboarding flow. setup_mode.py drops this sentinel file once nmcli has
# successfully joined the home Wi-Fi.
ConditionPathExists=/var/lib/bor-systems/wifi-configured

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
EnvironmentFile=/etc/bor-hanger.env
ExecStart=__PYTHON__ __HANGER_PY__
Restart=on-failure
RestartSec=10
# Allow GPIO access
SupplementaryGroups=gpio

[Install]
WantedBy=multi-user.target
